# udev/60-usb-no-suspend.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="1d6b", ATTRS{idProduct}=="0002", RUN+="/opt/bin/usb-no-suspend.sh"
